# docker/rootfs/etc/services.d/config_monitor/run
#!/usr/bin/with-contenv bash

umask ${UMASK}

cd ${WORKDIR}

exec \
    s6-setuidgid nexus .venv/bin/python config_monitor.py
